# dalla-stats-backup.service
[Unit]
Description=Export dalla-stats database

[Service]
Type=oneshot
ExecStart=/usr/local/dalla-stats2/backup.sh
